# derby.service
[Unit]
Description=Pinewood Derby
After=network.target

[Service]
Type=simple
User=nobody
ExecStart=derby

[Install]
WantedBy=multi-user.target
